# scripts/com.briefbot.dashboard-frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefbot.dashboard-frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/.venv/bin/python</string>
        <string>-m</string>
        <string>dashboard.backend.static_server</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>59000</string>
        <string>--dir</string>
        <string>__PROJECT_DIR__/dashboard/dist</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if the process exits -->
    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/data/logs/dashboard-frontend.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/data/logs/dashboard-frontend.log</string>
</dict>
</plist>
